# templogger.service
[Unit]
Description=Temperature Logger Script
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/path/to/your/script
ExecStart=/path/to/venv/bin/python /path/to/your/script/templogger.py
